# ansible/roles/backup_btrbk/templates/btrbk-job.timer.j2
[Unit]
Description=btrbk {{ job.description }} timer

[Timer]
OnCalendar={{ job.schedule }}
Persistent=true
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
